# deploy/services.d/app/run
#!/usr/bin/with-contenv sh
cd /usr/local/app
scrapyd
